# docker/rabbitmq/supervisor/conf.d/consumer_apm.conf
[program:consumer-apm]
command=php /var/www/floraqueen/bin/console rabbitmq:consumer apm
user=www-data 
autorestart=true
stopasgroup=true
stopsignal=QUIT
