# docker/rabbitmq/supervisor/conf.d/consumer_fuel_command.conf
[program:consumer-fuel-command]
command=php /var/www/floraqueen/bin/console rabbitmq:consumer fuel_command 
user=www-data 
autorestart=true
stopasgroup=true
stopsignal=QUIT
